# systemd/max-bridge.service
[Unit]
Description=MAX-Hermes Bridge - MAX Bot API to Hermes Agent bridge
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/opt/max-hermes
Environment=PATH=/opt/max-hermes/venv/bin:/usr/local/bin:/usr/bin
EnvironmentFile=/etc/max-bridge/.env
ExecStart=/opt/max-hermes/venv/bin/python -m src.main
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=max-bridge

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/max-bridge /opt/max-hermes /home/%i/.hermes/logs /home/%i/.hermes/sessions /home/%i/.hermes/workspace
PrivateTmp=true

[Install]
WantedBy=multi-user.target
